# plist/cn.elvea.redis.plist
<?xml version='1.0' encoding='UTF-8'?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version='1.0'>
    <dict>
        <key>Label</key>
        <string>cn.elvea.redis</string>
        <key>ProgramArguments</key>
        <array>
            <string>/Users/elvea/Tools/redis/redis/bin/redis-server</string>
            <string>/Users/elvea/Tools/redis/redis/redis.conf</string>
        </array>
        <key>RunAtLoad</key>
        <true/>
        <key>KeepAlive</key>
        <true/>
        <key>WorkingDirectory</key>
        <string>/Users/elvea/Tools/redis/redis</string>
        <key>StandardErrorPath</key>
        <string>/Users/elvea/Tools/redis/log/output.log</string>
        <key>StandardOutPath</key>
        <string>/Users/elvea/Tools/redis/log/output.log</string>
        <key>HardResourceLimits</key>
        <dict>
            <key>NumberOfFiles</key>
            <integer>1024</integer>
        </dict>
        <key>SoftResourceLimits</key>
        <dict>
            <key>NumberOfFiles</key>
            <integer>1024</integer>
        </dict>
    </dict>
</plist>
        s